# roles/systemd/templates/monitor_test.service.j2
[Unit]
Description=Monitor test process

[Service]
Type=oneshot
EnvironmentFile=-/etc/monitor_test.env
ExecStart=/usr/local/bin/monitor_test.sh